# contrib/com.agent-tui.daemon.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!--
  launchd plist for agent-tui daemon (macOS)

  Installation:
    1. Copy this file to ~/Library/LaunchAgents/
    2. Update the ProgramArguments path to your daemon location
    3. Load the service:
       launchctl load ~/Library/LaunchAgents/com.agent-tui.daemon.plist

  Management:
    launchctl list | grep agent-tui
    launchctl start com.agent-tui.daemon
    launchctl stop com.agent-tui.daemon
    launchctl unload ~/Library/LaunchAgents/com.agent-tui.daemon.plist
-->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.agent-tui.daemon</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/node</string>
        <string>/usr/local/share/agent-tui/daemon/dist/index.js</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENT_TUI_LOG_LEVEL</key>
        <string>info</string>
        <key>AGENT_TUI_TRANSPORT</key>
        <string>unix</string>
        <!-- Uncomment for TCP transport
        <key>AGENT_TUI_TRANSPORT</key>
        <string>tcp</string>
        <key>AGENT_TUI_TCP_PORT</key>
        <string>19847</string>
        -->
    </dict>

    <key>RunAtLoad</key>
    <false/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>5</integer>

    <key>StandardOutPath</key>
    <string>/tmp/agent-tui-daemon.stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/tmp/agent-tui-daemon.stderr.log</string>

    <key>ProcessType</key>
    <string>Background</string>

    <key>LowPriorityIO</key>
    <true/>

    <key>Nice</key>
    <integer>10</integer>
</dict>
</plist>
